# data/export/upstart_user/unicorn.conf.erb
start on starting <%= app %>
stop on stopping <%= app %>
respawn

chdir <%= engine.root %>
env HOME=<%= ENV['HOME'] %>
env WEB_CONCURRENCY=<%= num %>

script
  . $HOME/.bashrc
  <%= process.command %> >> <%= log %>/<%=name%>.log 2>&1
end script
