# data/export/upstart_user/unicorn.conf.erb
start on starting <%= app %>
stop on stopping <%= app %>
respawn

chdir <%= engine.root %>

exec sudo -i -u inspire bash -c "WEB_CONCURRENCY=<%= num %> bin/unicorn -c unicorn.rb >> <%= log %>/<%=name%>.log 2>&1"
